# config/systemd/user/rss.timer
[Unit]
Description=Run rss hourly

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
